# 99-hotplugged.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/bin/systemctl start nvidia_hotplug.service"
