# 99-hotplugged.rules
ACTION=="change", SUBSYSTEM="drm", RUN+="path/to/script.py"